# infra/systemd/industrial-lab-full@.service
[Unit]
Description=Manual bounded industrial research cycle (%i minutes)
After=network-online.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/industrial-risk-control
ExecStart=/home/ubuntu/industrial-risk-control/.venv/bin/python -m industrial_research_lab.cli --profile full --max-minutes %i --output artifacts/server-full-%i
