# deployment/systemd/industrial-risk-control-validation.service
[Unit]
Description=Industrial risk control validation
After=network-online.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/industrial-risk-control
ExecStart=/home/ubuntu/industrial-risk-control/.venv/bin/python -m pytest -q
